# rc/iked.in
#! /bin/sh
# $Id: rc.d-iked.in,v 1.3 2005/10/28 06:33:16 kamada Exp $
# iked rc.d script for NetBSD
#

# PROVIDE: ike
# REQUIRE: isdnd kdc ppp
# BEFORE:  SERVERS

. /etc/rc.subr

exec_prefix="@prefix@"
prefix="@prefix@"
name="iked"
rcvar=$name
command="@sbindir@/${name}"
extra_commands="reload"
required_files="@sysconfdir@/iked.conf"
reload_cmd="@sbindir@/ikectl ${extra_commands}"

load_rc_config $name
run_rc_command "$1"
