# rc/iked.in
#! /bin/sh
# $Id: rc.d-iked.in,v 1.3 2005/10/28 06:33:16 kamada Exp $
# iked rc.d script for NetBSD
#

# PROVIDE: ike
# REQUIRE: isdnd kdc ppp
# BEFORE:  SERVERS

. /etc/rc.subr

# We only need to stop the process owned by root
check_process_root()
{
	_procname=$1
	_interpreter=$2
	if [ -z "$_procname" ]; then
		err 3 'USAGE: check_process procname [interpreter]'
	fi
	_find_processes $_procname ${_interpreter:-.} '-Uxroot'
}

exec_prefix="@prefix@"
prefix="@prefix@"
name="iked"
rcvar=$name
command="@sbindir@/${name}"
check_process_root ${name}
extra_commands="reload"
required_files="@sysconfdir@/iked.conf"
reload_cmd="@sbindir@/ikectl ${extra_commands}"

load_rc_config $name
run_rc_command "$1"
